# source/services/minecraft-prod.service
[Unit]
Description=Minecraft Server Service

[Service]
Type=simple
Restart=always
RestartSec=3
User=minecraftd
ExecStart=/home/minecraftd/server/bin/initialize-prod

[Install]
WantedBy=multi-user.target
